# rubifo.service
[Unit]
Description=Rubifo - Rubika Auto-Forward Bot
Documentation=https://rubifo.ir/docs
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=rubifo
Group=rubifo
WorkingDirectory=/opt/rubifo

# Environment variables
EnvironmentFile=/opt/rubifo/.env
Environment="PATH=/opt/rubifo/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

# Start command
ExecStart=/opt/rubifo/venv/bin/python -m src.bot.main

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Security and logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rubifo

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Process management
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

# Auto-start on boot
Alias=rubifo.service
